# script/goSync
#! /bin/sh
# /etc/init.d/goSync

### BEGIN INIT INFO
# Provides:          goSync
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Simple script to start a program at boot
# Description:       A simple script from www.stuffaboutcode.com which will start / stop a program a boot / shutdown.
### END INIT INFO

# If you want a command to always run, put it here

# Carry out specific functions when asked to by the system
case "$1" in
  start)
    echo "Starting syncDrive"
    sudo python /home/pi/Documents/techlab_iot/script/syncDrive.py  >/tmp/syncLog_error.txt 2>/tmp/syncLog_normal.txt &
    ;;
  stop)
    echo "Stopping syncDrive"
    sudo kill -9 $(ps aux | grep '[s]yncDrive.py' | awk '{print $2}')
    sudo mv /tmp/syncLog_normal.txt /tmp/syncLog_normal.bak
    sudo mv /tmp/syncLog_error.txt /tmp/syncLog_error.bak
    ;;
  *)
    echo "Usage: /etc/init.d/goSync {start|stop}"
    exit 1
    ;;
esac

exit 0
